# se8200.conf
[program:se8200]
command=/usr/bin/sh /opt/run-se8200.sh
redirect_stderr=true
autorestart=unexpected
stdout_logfile=/var/log/se8200.log
