# sshd.conf
[program:sshd]
command=/usr/sbin/sshd -Df /etc/ssh/sshd_config
redirect_stderr=true
autorestart=unexpected
stdout_logfile=/var/log/sshd.log
